# docker/s6-overlay/s6-rc.d/mekong-gateway/run
#!/bin/sh
cd /workspace
exec s6-setuidgid claude env \
  HOME=/home/claude \
  MEKONG_HOME=/workspace \
  DISPLAY=:99 \
  python3 -m uvicorn src.core.gateway:app --host 0.0.0.0 --port 8000
