# scripts/com.mekong.9router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mekong.9router</string>

    <key>ProgramArguments</key>
    <array>
        <string>__9ROUTER_PATH__</string>
        <string>--no-browser</string>
        <string>--skip-update</string>
        <string>--port</string>
        <string>20128</string>
        <string>--log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/9router_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/9router_stderr.log</string>
</dict>
</plist>
