# scripts/com.antigravity.macook-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.macook-guardian</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/macbookprom1/mekong-cli/scripts/macook-guardian.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/macbookprom1/mekong-cli/logs/macook-guardian.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookprom1/mekong-cli/logs/macook-guardian.stderr.log</string>
</dict>
</plist>
