# scripts/com.mekong.tom-hum-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mekong.tom-hum-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__MEKONG_DIR__/apps/openclaw-worker/task-watcher.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__MEKONG_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_BASE_URL</key>
        <string>http://127.0.0.1:8080</string>
        <key>MEKONG_DIR</key>
        <string>__MEKONG_DIR__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/tom_hum_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/tom_hum_stderr.log</string>
</dict>
</plist>
